# .config/systemd/user/rclone-sync-drive.path
[Unit]
Description=Dropbox sync monitor

[Path]
PathChanged=%h/cloud
PathModified=%h/cloud
Unit=rclone-sync-drive.service

[Install]
WantedBy=multi-user.target
