# .config/systemd/user/rclone-dropbox.path
[Unit]
Description=Dropbox sync monitor

[Path]
PathChanged=%h/cloud
PathModified=%h/cloud
Unit=rclone-dropbox.service

[Install]
WantedBy=multi-user.target
